# jjh/com.tealeaf.jjh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tealeaf.jjh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tealeaf/Local_Code/GitHub/launchd/jjh/run.sh</string>
    </array>

    <!-- 매일 KST 09:00 실행 (시스템 시간대: Asia/Seoul) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/tealeaf/Local_Code/GitHub/launchd/jjh/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tealeaf/Local_Code/GitHub/launchd/jjh/logs/stderr.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/tealeaf/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/tealeaf</string>
        <key>LANG</key>
        <string>ko_KR.UTF-8</string>
    </dict>
</dict>
</plist>
